# systemd/ai-link-scraper-daemon.service
[Unit]
Description=AI Link Scraper Daemon
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/ai-link-scraper
Environment=PATH=/path/to/your/python/bin
ExecStart=/path/to/your/python/bin/python daemon.py --interval 60
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables (alternatively, use EnvironmentFile)
# Environment="SLACK_BOT_TOKEN=your-token-here"
# Environment="OPENAI_API_KEY=your-key-here"
# Or use:
# EnvironmentFile=/path/to/ai-link-scraper/.env

[Install]
WantedBy=multi-user.target
